# templates/default/etcd.service.erb
[Unit]
Description=etcd

[Service]
<% @environment.each do |key, value| -%>
Environment=<%= key %>=<%= value %>
<% end %>
ExecStart=/usr/local/bin/etcd

[Install]
WantedBy=multi-user.target
